# fs/etc/services.d/php/run
#!/usr/bin/with-contenv sh

exec php-fpm5 -F;